# sandbox/services/lss-sync/run
#!/usr/bin/with-contenv bash
# lss-sync — File-watcher daemon that keeps the semantic search index in sync.
#
# Uses watchdog (inotify/FSEvents) to detect file changes in real time and
# triggers re-indexing with debounce. Changes are picked up within seconds.
#
# Watches:
#   /workspace       — User files and projects
#   /workspace/.kortix — Agent memory

exec s6-setuidgid abc \
  env HOME=/workspace \
      OPENAI_API_KEY="${OPENAI_API_KEY}" \
      LSS_DIR=/workspace/.lss \
      PATH="/lsiopy/bin:/usr/local/bin:/usr/bin:/bin" \
  lss-sync \
    --watch /workspace \
    --exclude node_modules \
    --exclude .git \
    --exclude __pycache__ \
    --startup-delay 15 \
    --debounce 2
