# s6-inspect/s6-rc.d/svc-static-web/finish
#!/bin/bash
# Kill any stale bun process holding port 3211 before restart
fuser -k 3211/tcp 2>/dev/null || true
sleep 2
